# deploy/systemd/oxidity-searcher.service
[Unit]
Description=Oxidity Searcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=john
Group=john
WorkingDirectory=/opt/oxidity/bin
EnvironmentFile=-/etc/oxidity/searcher.env
Environment=RUST_LOG=info
ExecStart=/opt/oxidity/bin/oxidity-searcher
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
