# deploy/systemd/oxidity-wallet-backend.service
[Unit]
Description=Oxidity Wallet Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=john
Group=john
WorkingDirectory=/opt/oxidity/bin
EnvironmentFile=-/etc/oxidity/wallet-backend.env
Environment=RUST_LOG=info
ExecStart=/opt/oxidity/bin/oxidity_wallet_backend
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
