# deploy/systemd/wallet-service.service
[Unit]
Description=Oxidity Wallet Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=john
WorkingDirectory=/home/john/oxidity-searcher
EnvironmentFile=-/home/john/oxidity-searcher/.env
Environment=RUST_LOG=info
ExecStart=/home/john/oxidity-searcher/target/release/wallet_service --bind 127.0.0.1 --port 9555
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
